# packaging/core/postinst
#!/bin/bash
# Project: Crankshaft
# This file is part of Crankshaft project.
# Copyright (C) 2025 OpenCarDev Team
#
#  Crankshaft is free software: you can redistribute it and/or modify
#  it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 3 of the License, or
#  (at your option) any later version.
#
#  Crankshaft is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU General Public License for more details.
#
#  You should have received a copy of the GNU General Public License
#  along with Crankshaft. If not, see <http://www.gnu.org/licenses/>.

set -e

case "$1" in
    configure)
        # Create crankshaft user if it doesn't exist
        if ! id -u crankshaft > /dev/null 2>&1; then
            useradd --system --no-create-home --shell /bin/false crankshaft
        fi
        
        # Create plugdev group if it doesn't exist (for USB device access)
        if ! getent group plugdev > /dev/null 2>&1; then
            groupadd --system plugdev
        fi
        
        # Add crankshaft user to plugdev group for USB access
        usermod -a -G plugdev crankshaft 2>/dev/null || true
        
        # Set proper ownership for config directory
        chown -R crankshaft:crankshaft /etc/crankshaft
        chmod 755 /etc/crankshaft
        chmod 644 /etc/crankshaft/crankshaft.json
        
        # Create runtime directories
        mkdir -p /var/lib/crankshaft /var/log/crankshaft
        chown -R crankshaft:crankshaft /var/lib/crankshaft /var/log/crankshaft
        
        # Ensure USB action script is executable
        chmod +x /usr/lib/crankshaft/usb_action.sh
        
        # Reload udev rules
        udevadm control --reload-rules 2>/dev/null || true
        udevadm trigger 2>/dev/null || true
        
        # Reload systemd and enable service when running under systemd
        if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable crankshaft-core.service
        else
            echo "Systemd not detected; skipping service enable." >&2
        fi
        
        echo "Crankshaft Core installed successfully"
        echo "Service enabled and will start on boot"
        echo "USB device detection configured"
        echo "To start now: sudo systemctl start crankshaft-core"
    ;;
esac

exit 0
