# packaging/core/prerm
#!/bin/sh
# Project: Crankshaft
# This file is part of Crankshaft project.
# Copyright (C) 2025 OpenCarDev Team
#
#  Crankshaft is free software: you can redistribute it and/or modify
#  it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 3 of the License, or
#  (at your option) any later version.
#
#  Crankshaft is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU General Public License for more details.
#
#  You should have received a copy of the GNU General Public License
#  along with Crankshaft. If not, see <http://www.gnu.org/licenses/>.

set -e
export DEBIAN_FRONTEND=noninteractive

case "$1" in
    remove|deconfigure)
        # Stop and disable systemd service (if systemd is present)
        if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
            systemctl stop crankshaft-core.service >/dev/null 2>&1 || true
            systemctl disable crankshaft-core.service >/dev/null 2>&1 || true
        fi

        # Kill any remaining processes (best effort)
        pkill -f '^/usr/bin/crankshaft-core' >/dev/null 2>&1 || true
        pkill -f '[c]rankshaft-core' >/dev/null 2>&1 || true

        echo "Crankshaft Core pre-removal cleanup done"
    ;;
    upgrade)
        # During upgrade, only stop the service gracefully
        if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
            systemctl stop crankshaft-core.service >/dev/null 2>&1 || true
        fi
        
        # Wait for process to exit (with timeout)
        for i in 1 2 3 4 5; do
            if ! pgrep -f '[c]rankshaft-core' >/dev/null 2>&1; then
                break
            fi
            sleep 1
        done
        
        # Force kill if still running
        pkill -9 -f '[c]rankshaft-core' >/dev/null 2>&1 || true
        
        echo "Crankshaft Core stopped for upgrade"
    ;;
esac

exit 0
