# packaging/ui/crankshaft-ui.service
# Project: Crankshaft
# This file is part of Crankshaft project.
# Copyright (C) 2025 OpenCarDev Team
#
#  Crankshaft is free software: you can redistribute it and/or modify
#  it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 3 of the License, or
#  (at your option) any later version.
#
#  Crankshaft is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU General Public License for more details.
#
#  You should have received a copy of the GNU General Public License
#  along with Crankshaft. If not, see <http://www.gnu.org/licenses/>.

[Unit]
Description=Crankshaft UI Service
Documentation=https://github.com/opencardev/crankshaft.core
After=crankshaft-core.service graphical.target
Requires=crankshaft-core.service
Wants=graphical.target

[Service]
Type=simple
User=crankshaft
Group=crankshaft
ExecStart=/usr/bin/crankshaft-ui
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment for Qt platform
Environment="QT_QPA_PLATFORM=eglfs"
Environment="QT_QPA_EGLFS_INTEGRATION=eglfs_kms"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/crankshaft /var/log/crankshaft
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=graphical.target
